# etc/xdg/autostart/lxqt-klipper.desktop
[Desktop Entry]
Type=Application
Name=Klipper
Exec=klipper
OnlyShowIn=LXQt;
